# src/ai/backend/install/pyinfra/deploy/cores/appproxy/templates/systemd/backendai-appproxy-delete-circuits.timer.j2
[Unit]
Description=Backend.AI AppProxy Delete Unused Circuits Timer (Every 5 minutes)

[Timer]
OnCalendar=*:0/5:00
Persistent=true
Unit=backendai-appproxy-delete-circuits.service

[Install]
WantedBy=timers.target
