# src/ai/backend/install/pyinfra/deploy/cores/appproxy/templates/systemd/backendai-appproxy-cleanup-traefik.timer.j2
[Unit]
Description=Backend.AI AppProxy Cleanup Traefik Routes Timer (Every 2 minutes)

[Timer]
OnCalendar=*:0/2:00
Persistent=true
Unit=backendai-appproxy-cleanup-traefik.service

[Install]
WantedBy=timers.target
