# src/ai/backend/install/pyinfra/deploy/cores/manager/templates/systemd/backendai-manager-restart.timer.j2
[Unit]
Description=Backend.AI Manager Restart Timer (Daily at {{ restart_time }})

[Timer]
OnCalendar=*-*-* {{ restart_time }}
Persistent=true
Unit=backendai-manager-restart.service

[Install]
WantedBy=timers.target
