# src/ai/backend/install/pyinfra/deploy/cores/manager/templates/systemd/backendai-manager-restart.service.j2
[Unit]
Description=Backend.AI Manager Restart Service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart backendai-manager.service
User=root
