# src/ai/backend/install/pyinfra/deploy/cores/manager/templates/systemd/backendai-manager-data-backup.service.j2
[Unit]
Description=Backend.AI Manager Data Backup Service

[Service]
Type=oneshot
ExecStart={{ home_dir }}/scripts/backup_bai_data.sh
User={{ bai_user }}
WorkingDirectory={{ home_dir }}
StandardOutput=journal
StandardError=journal
